# rfcomm-bt.conf
[program:rfcomm-bt]
command = /usr/bin/python /home/pi/rfcomm-server.py
autostart=true
